# ytmusic-web.service
[Unit]
Description=Espacoviva Infinite Music Server (v4.0 DevOps)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=bruno
Group=bruno
WorkingDirectory=/home/bruno/espacoviva-deploy
Environment=HOME=/home/bruno
Environment=PYTHONUNBUFFERED=1

# Pre-execution cleanup
ExecStartPre=-/usr/bin/killall -9 mpv
ExecStartPre=-/usr/bin/rm -f /tmp/mpv-ytmusic

# Main process
ExecStart=/usr/bin/python3 /usr/local/bin/ytmusic-web.py

# Failure handling & Resilience
Restart=always
RestartSec=5s
StartLimitIntervalSec=0

# Resource priority (Industrial grade)
Nice=-10
OOMScoreAdjust=-1000
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
